# python/crontab
0 17 * * * root python3 start.py >> /var/log/$(date +%Y_%m_%d).log 2>&1
